# centos-local/etc/systemd/system/hugepage-fix.service
[Unit]
Before=vpp.service

[Service]
ExecStart=/root/setup.sh

[Install]
WantedBy=default.target
