# resources/tools/disk-image-builder/csit-sut-dcr/ubuntu_1804_amd64/supervisord.conf
[unix_http_server]
file  = /run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///run/supervisor.sock

[supervisord]
pidfile = /run/supervisord.pid
identifier = supervisor
directory = /run
logfile=/var/log/supervisord.log
loglevel=debug
nodaemon=false

[program:vpp]
command=/usr/bin/vpp -c /etc/vpp/startup.conf
autorestart=true
redirect_stderr=true
priority=1
