# resources/tools/testbed-setup/ansible/roles/cobbler/files/supervisord/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock

[supervisord]
pidfile=/var/run/supervisord.pid
identifier=supervisor
directory=/run
logfile=/var/log/supervisord.log
loglevel=debug
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:cobblerd]
command=/usr/bin/cobblerd -F
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log

[program:tftpd]
command=/usr/sbin/in.tftpd --foreground --verbose --user root --permissive --blocksize 1380 --address 0.0.0.0:69 --secure /var/lib/tftpboot
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log

[program:rsyncd]
command=/usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf
stopsignal=QUIT
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log

[program:dhcpd]
command=/usr/sbin/dhcpd -f -cf /etc/dhcp/dhcpd.conf -user dhcpd -group dhcpd --no-pid
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log
